# systemd/qa-cell-edge-agent.target
[Unit]
Description=QA Cell Edge Agent (all processes)
After=network-online.target
Wants=qa-cell-edge-agent.service

[Install]
WantedBy=multi-user.target
