# systemd/qa-defect-detection.service
[Unit]
Description=QA Cell - Defect Detection (Process 2)
After=network-online.target qa-sensor-push.service
Requires=network-online.target
PartOf=qa-cell-edge-agent.target

[Service]
Type=simple
User=jetson
WorkingDirectory=/opt/qa-cell-edge-agent
EnvironmentFile=/opt/qa-cell-edge-agent/.env
ExecStart=/opt/qa-cell-edge-agent/venv/bin/python -m qa_cell_edge_agent.processes.defect_detection
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qa-defect-detection

[Install]
WantedBy=qa-cell-edge-agent.target
