# systemd/qa-model-upgrade.service
[Unit]
Description=QA Cell - Model Upgrade (Process 3)
After=network-online.target
Requires=network-online.target
PartOf=qa-cell-edge-agent.target

[Service]
Type=simple
User=jetson
WorkingDirectory=/opt/qa-cell-edge-agent
EnvironmentFile=/opt/qa-cell-edge-agent/.env
ExecStart=/opt/qa-cell-edge-agent/venv/bin/python -m qa_cell_edge_agent.processes.model_upgrade
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=qa-model-upgrade

[Install]
WantedBy=qa-cell-edge-agent.target
